# crontab
*/1 * * * * python3 /cron.py >> /logfile